# deploy/csops-roadmap.service
[Unit]
Description=CS Ops Roadmap Dashboard (Flask)
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/csops-roadmap
ExecStart=/usr/bin/python3 /var/www/csops-roadmap/api/app.py
Restart=always
RestartSec=3
Environment=FLASK_ENV=production
StandardOutput=journal
StandardError=journal
SyslogIdentifier=csops-roadmap

[Install]
WantedBy=multi-user.target
